# board/dc-rpi3/rootfs_overlay/etc/cron.d/system
# root crontab
# DO NOT EDIT THIS FILE MANUALLY! USE crontab -e INSTEAD

# man 1 crontab for acceptable formats:
#    <minute> <hour> <day> <month> <dow> <tags and command>
#    <@freq> <tags and command>

# SYSTEM DAILY/WEEKLY/... FOLDERS
@hourly         ID=sys-hourly   /bin/run-parts /etc/cron.hourly
@daily          ID=sys-daily    /bin/run-parts /etc/cron.daily
@weekly         ID=sys-weekly   /bin/run-parts /etc/cron.weekly
@monthly        ID=sys-monthly  /bin/run-parts /etc/cron.monthly

*/1 * * * * /usr/bin/php-cgi /usr/bin/wifiAutoReconnect.php

